# fixtures/glass-desktop/glass-desktop.service
[Unit]
Description=Glass XFCE desktop
DefaultDependencies=no
After=ccx3-stage2.service
Wants=ccx3-stage2.service

[Service]
Type=simple
ExecStart=/usr/local/sbin/glass-session
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=glass.target
